# deploy/systemd/vignesh-agent.service
[Unit]
Description=Vignesh WhatsApp Agent
After=network-online.target wabot.service
Wants=network-online.target

[Service]
Type=simple
User=vignesh
Group=vignesh
WorkingDirectory=/opt/vignesh-agent
EnvironmentFile=/opt/vignesh-agent/.env
Environment=PYTHONUNBUFFERED=1
Environment=UV_CACHE_DIR=/opt/vignesh-agent/.uv-cache
ExecStart=/usr/local/bin/uv run --frozen --no-sync python main.py
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/vignesh-agent/data /opt/vignesh-agent/.uv-cache

[Install]
WantedBy=multi-user.target
